# archnas/fs/etc/profile.d/userenv.sh
# Variables for the user's environment and user experience
export PAGER=ov
export EDITOR=nvim

